# init/smack.mount
[Unit]
Description=Smack filesystem mounting
WantedBy=local-fs.target

[Mount]
What=smackfs
Where=/smack
Type=smackfs
