# init/smack.mount
[Unit]
Description=Smack filesystem mounting

[Mount]
What=smackfs
Where=/smack
Type=smackfs
